# src/supervisor/conf.d/papercut.conf
[program:papercut]
command=/etc/init.d/papercut console
priority=900
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=root
autorestart=true